# resources/supervisor/laravel-worker.conf
[program:laravel-worker-[connection]-[queue]]
process_name=%(program_name)s_%(process_num)02d
command=php [application_root]/artisan queue:work [connection] --queue=[queue] --sleep=5 --timeout=7200 --tries=1
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=[application_root]/storage/logs/worker.log
stopwaitsecs=3600
